# alpine/ninjatorrent.ini
[supervisord]
nodaemon=true

[program:deluged]
command=deluged -d

[program:deluge-web]
command=deluge-web -d

#[program:flexget]
#command=flexget daemon start
